# deploy/systemd/morningnews-youtube.service
[Unit]
Description=MorningNews YouTube daily follow-up job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=REPLACE_WITH_RUN_USER
WorkingDirectory=REPLACE_WITH_MORNINGNEWS_YOUTUBE_DIR
EnvironmentFile=REPLACE_WITH_MORNINGNEWS_YOUTUBE_DIR/.env
ExecStart=/bin/bash REPLACE_WITH_MORNINGNEWS_YOUTUBE_DIR/bin/run_daily.sh

[Install]
WantedBy=multi-user.target