# deploy/systemd/morningnews.service
[Unit]
Description=MorningNews daily generation job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=glamboy
WorkingDirectory=/data/projects/morningnews
EnvironmentFile=/data/projects/morningnews/.env
ExecStart=/bin/bash /data/projects/morningnews/run_morningnews.sh

[Install]
WantedBy=multi-user.target