# deploy/systemd/morningnews-youtube.service
[Unit]
Description=MorningNews YouTube daily follow-up job
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=REPLACE_WITH_RUN_USER
WorkingDirectory=/data/projects/morningnews-youtube
EnvironmentFile=/data/projects/morningnews-youtube/.env
ExecStart=/bin/bash /data/projects/morningnews-youtube/bin/run_daily.sh

[Install]
WantedBy=multi-user.target