# deploy/systemd/morningnews.timer
[Unit]
Description=Run MorningNews every day at 06:00 KST

[Timer]
OnCalendar=*-*-* 06:00:00 Asia/Seoul
Persistent=true
Unit=morningnews.service

[Install]
WantedBy=timers.target